# gpio_emergency/rootfs/etc/services.d/gpio-emergency/run
#!/usr/bin/with-contenv bashio
set -e
export MQTT_HOST="$(bashio::config 'mqtt_host')"
export GPIO_PIN="$(bashio::config 'gpio_pin')"
exec python3 /emergency_button.py
