# gpio_emergency/rootfs/etc/services.d/gpio-emergency/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting GPIO Emergency service"
exec python3 /emergency_button.py
